# systemd/agri-context-bridge.service
[Unit]
Description=EdgeLink RK3568 Agri Context Bridge
After=network-online.target openclaw-gateway.service
Wants=network-online.target

[Service]
Type=simple
User=linaro
Group=linaro
WorkingDirectory=/home/linaro/project/EdgeLink_RK3568
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 /home/linaro/project/EdgeLink_RK3568/src/agri_context_bridge.py --config /home/linaro/project/EdgeLink_RK3568/config/agri-context-bridge.ini
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
